# recovery-ui.init.in
#!/bin/sh
#
# Copyright (C) 2017 Dream Property GmbH, Germany
#                    https://dreambox.de/
#
### BEGIN INIT INFO
# Provides:          recovery_ui
# Required-Start:    $local_fs
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start recovery-ui at boot time
# Description:       Enable service provided by recovery-ui.
### END INIT INFO

test -x @bindir@/recovery-ui || exit 0

case "$1" in
  start)
    echo -n "Starting recovery-ui"
    start-stop-daemon -q -p @runstatedir@/recovery-ui.pid -x @bindir@/recovery-ui -S -b -m
    echo "."
    ;;
  stop)
    echo -n "Stopping recovery-ui"
    start-stop-daemon -q -p @runstatedir@/recovery-ui.pid -x @bindir@/recovery-ui -K
    echo "."
    ;;
  restart|reload|force-reload)
    $0 stop
    sleep 1
    $0 start
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|reload|force-reload}"
    exit 1
    ;;
esac

exit 0
